# trust-agent/mtwilson-tpm2-packages/src/resources/tss2/tcsd2.service
[Unit]
Description=tss2 resource manager
After=syslog.target

[Service]
ExecStart=/usr/local/sbin/resourcemgr $RESOURCED_OPTS

[Install]
WantedBy=multi-user.target
Alias=tss2resourced.service
